# debian-control/postrm.in
#!/usr/bin/env bash

set +e

systemctl stop ${VXLNETWORK_SERVICE} 2>/dev/null
systemctl disable ${VXLNETWORK_SERVICE} 2>/dev/null
rm -rf /etc/systemd/system/${VXLNETWORK_SERVICE}
systemctl daemon-reload

set -e
